# systemd/autopatchd.timer.in
[Unit]
Description=autopatchd scheduled updates

[Timer]
# Default to every Sunday 02:00, can be replaced in setup()
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
AccuracySec=1m
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target