# systemd/autopatchd.service.in
[Unit]
Description=autopatchd - automated system patch daemon
After=network-online.target

[Service]
Type=oneshot
LoadCredential=mailpass:/etc/autopatchd/smtp-password.cred
ExecStart=/usr/local/bin/autopatchd run
StandardOutput=journal
StandardError=journal